# udev.rules.d/unused/70-wifi-powersave.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlp3s0", RUN+="/usr/bin/iw dev %k set power_save on"
